# src/bash/crontab
*/3 * * * * /home/wilson/cryptohamster/run_hamsterwheel_handler.sh
30 * * * * > /home/wilson/log/hamsterwheel.log
0 13 * * 6 > /home/wilson/log/hamsterwheel_handler.log